# codeproject-ai-server-dev/rootfs/etc/s6-overlay/s6-rc.d/init-codeprojectai/run
#!/command/with-contenv bashio
# shellcheck shell=bash

mkdir -p "/config/runtimes"
mkdir -p "/config/modules"

# Check if the runtimes directory is a symlink, if not create it
if [ ! -L "/app/server/runtimes"]; then
    rm -rf "/app/server/runtimes"
    ln -s "/config/runtimes" "/app/server/runtimes"
fi

# Check if the modules directory is a symlink, if not create it
if [ ! -L "/app/server/modules"]; then
    rm -rf "/app/server/modules"
    ln -s "/config/modules" "/app/server/modules"
fi

if [ ! -L "/app/server/appsettings.json" ]; then  
    cp /app/server/appsettings.json /config/appsettings.json
    rm /app/server/appsettings.json
    ln -s /config/appsettings.json /app/server/appsettings.json
fi

cd /app/server || exit 1
exec bash ../setup.sh
