# codeprojectai/rootfs/etc/s6-overlay/s6-rc.d/init-codeprojectai/run
#!/command/with-contenv bashio
# shellcheck shell=bash

cd /app/server || exit 1
exec bash ../setup.sh
